# resources/systemd/hblock.service
[Unit]
Description=hBlock

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hblock
