# resources/systemd/hblock.service
[Unit]
Description=Update hosts file to blackhole ad sites

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hblock
