# resources/systemd/hblock.timer
[Unit]
Description=hBlock
After=network.target network-online.target
Requires=network-online.target

[Timer]
OnCalendar=*-*-* 00:00:00
RandomizedDelaySec=3600
Persistent=true
Unit=hblock.service

[Install]
WantedBy=timers.target
