# resources/systemd/hblock.service
[Unit]
Description=hBlock

[Service]
Type=oneshot
TimeoutSec=300
ExecStart=/bin/sh -c 'until ping -c1 github.com; do sleep 30; done && hblock'
